# sv/elasticsearch/run
#!/bin/bash

#sv start socklog-unix || exit 1

source /etc/envvars

exec 2>&1
chpst -u elasticsearch /elasticsearch/bin/elasticsearch
